# ub2204-pihole/files/pihole-backups.timer
[Unit]
Description=Backup pihole data

[Timer]
OnCalendar=17:30

[Install]
WantedBy=timers.target
