# ub2004-firefly3/files/firefly3recur.service
[Unit]
Description=Firefly III recurring transactions

[Service]
Type=oneshot
ExecStart=/usr/bin/php /var/www/firefly/artisan firefly-iii:cron -qn