# ub2004-firefly3/files/firefly3recur.timer
[Unit]
Description=Firefly III recurring transactions

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target